# etc/systemd/system/python312-default.service
[Unit]
Description=Set Python 3.12 as default
After=multi-user.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'source /etc/profile.d/python312_default.sh && update-alternatives --set python /usr/local/bin/python3.12 2>/dev/null || true && update-alternatives --set python3 /usr/local/bin/python3.12 2>/dev/null || true'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
